# systemd/lims-backup.timer
[Unit]
Description=LabScientific LIMS Backup Timer
Documentation=https://github.com/your-org/labscientific-lims
Requires=lims-backup.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target